# etc/systemd/user/scrcpy.service
[Service]
ExecStart=/home/pi/auto-scrcpy.sh
User=pi
StandardOutput=tty1
